# pi/systemd/vigil-recorder.service
[Unit]
Description=vigil recorder — capture and segment webcam video
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=vigil
Group=vigil
WorkingDirectory=/opt/vigil
ExecStart=/opt/vigil/scripts/start_recording.sh --config /etc/vigil/vigil.yaml
Restart=on-failure
RestartSec=5
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
